# cron/egonsystem.cron
# Every 10 minutes
*/10 * * * * /usr/local/bin/python /app/scripts/collect_snapshot.py >> /var/log/egonsystem-cron.log 2>&1

